# external/daemon_script/apollo-lircd
#!/bin/bash
### BEGIN INIT INFO
# Provides: apollo-lircd
# Required-Start:   $remote_fs $syslog $time
# Required-Stop:    $remote_fs $syslog $time
# Default-Start:    2 3 4 5
# Default-Stop:     0 1 6
# Short-Description: APOLLO-lircd
### END INIT INFO

do_start()
{
    serial_list=`dmesg | grep "Product: FT230X Basic UART" -A 2 | grep "SerialNumber: " | awk -F ' ' '{print $6}'`

    index=0
    for i in $serial_list; do
        sudo lircd -Hftdix -dserial=$i -o /usr/local/var/run/lirc/lircd$index -P /usr/local/var/run/lirc/lircd$index.pid
        index=`expr $index + 1`
    done
}

do_stop()
{
    for(( i=0; i<=4; i++ )); do
        if [ -f /usr/local/var/run/lirc/lircd$i.pid ]; then
            pkill -F /usr/local/var/run/lirc/lircd$i.pid
        fi
    done
}

#########################################################################################################

case "$1" in
start)
    do_start
    ;;

stop)
    do_stop
    ;;

restart)
    do_stop
    sleep 1
    do_start
    ;;

*)

echo "Usage : service apollo-lircd {start|stop|restart}" >&2

esac
exit 0
